# conf/crontab
CRON_TZ=Asia/Tokyo
*/1 * * * * /usr/src/cron/exec_job.sh
