# conf/crontab
CRON_TZ=Asia/Tokyo
0 9,12,15,18,21 * * * /usr/src/cron/exec_job.sh
